# cron.d/crontab
*/1 * * * * BASH_ENV=/temp-ram-disk/container.env run-parts --report --exit-on-error --verbose --regex '.*.sh$' /opt/galaxy/cron.d/cron.minute > /var/log/cron.log 2>&1
0 */1 * * * BASH_ENV=/temp-ram-disk/container.env run-parts --report --exit-on-error --verbose --regex '.*.sh$' /opt/galaxy/cron.d/cron.hourly > /var/log/cron.log 2>&1
0 2 */1 * * BASH_ENV=/temp-ram-disk/container.env run-parts --report --exit-on-error --verbose --regex '.*.sh$' /opt/galaxy/cron.d/cron.daily > //var/log/cron.log 2>&1
@reboot BASH_ENV=/temp-ram-disk/container.env run-parts --report --exit-on-error --verbose --regex '.*.sh$' /opt/galaxy/cron.d/cron.once > /var/log/cron.log 2>&1
